# etc/update-motd.d/50-amazon-linux-extras-news
#!/bin/bash
timeout 5s amazon-linux-extras system_motd 2>/dev/null
